# overlay/etc/systemd/system/multi-user.target.wants/hostapd.service
[Unit]
Description=hostapd
After=wlan1.service

[Service]
RemainAfterExit=yes
ExecStart=/usr/sbin/hostapd -B /etc/hostapd.conf

[Install]
WantedBy=multi-user.target