# overlay/etc/systemd/system/multi-user.target.wants/snoopy.service
[Unit]
Description=snoopy service
Wants=network-manager.service
After=network-manager.service
ConditionPathExists=/etc/snoopy.conf

[Service]
ExecStart=/var/snoopy/snoopy -config=/etc/snoopy.conf

[Install]
WantedBy=multi-user.target